# router_firmware_download/xiaomi/outpath1/outpath/miwifi_ra70_firmware_cc424_1.0.168/debug/fs/etc/hotplug.d/button/02-mesh
#!/bin/sh

if [ "$ACTION" = "released" -a "$BUTTON" = "mesh" ]; then
	ubus call miwifi-discovery start_mesh
fi
