# router_firmware_download/xiaomi/outpath1/miwifi_ra70_firmware_cc424_1.0.168/debug/fs/etc/init.d/miled
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=98
NAME=miled
USE_PROCD=1
PROG=/usr/sbin/miled


start_service(){
    [ -f "$PROG" ] && {
	procd_open_instance
	procd_set_param command ${PROG}
	procd_set_param respawn
	procd_close_instance
    }
}

stop_service() {
    service_stop ${PROG}
}

restart_service() {
    service_reload ${PROG}
}
